# doc/log21/plc/qplc.service
[Unit]
Description=QPLC Startup Script
After=network.target

[Service]
ExecStart=/opt/sorter/plc/start_qplc.sh

[Install]
WantedBy=default.target
